# integrations/mc/systemd/mc-node.service
[Unit]
Description=MissionControl Node Agent
Documentation=https://github.com/missioncontrol-ai/missioncontrol
Wants=network-online.target
After=network-online.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/default/%N
EnvironmentFile=-/etc/sysconfig/%N
EnvironmentFile=-/etc/missioncontrol/%N.env
ExecStart=/usr/local/bin/mc node run
Restart=always
RestartSec=5s
KillMode=control-group
TimeoutStartSec=0
LimitNOFILE=1048576
